# bluetooth-autoconnector.service
[Unit]
Description=connect to device 
 
[Service]
ExecStart=/usr/local/bin/bluetooth-autoconnector.ps1
